# systemd/redshift-login.timer
[Timer]
Unit=redshift-timestamped.service
#Executed once after login
OnActiveSec=20
Persistent=false

[Install]
WantedBy=default.target

